# oracle-bot/deploy/nobd-oracle.service
[Unit]
Description=NOBD MvC2 Oracle — Discord bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/nobd-oracle/oracle-bot/discord-bot
# Secrets live here on the VPS only (chmod 600) — NEVER in git. See deploy/.env.example.
EnvironmentFile=/opt/nobd-oracle/oracle-bot/discord-bot/.env
ExecStart=/opt/nobd-oracle/venv/bin/python /opt/nobd-oracle/oracle-bot/discord-bot/bot.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
